# riemann-puppet-ping.upstart
# riemann-puppet-ping
description "riemann puppet icmp health check"
start on runlevel [235]
stop on runlevel [016]
respawn
exec riemann-puppet-ping
